# com.example.mac2mqttd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.mac2mqttd</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/mac2mqtt/mac2mqttd</string>
        <string>--config</string>
        <string>/Users/USERNAME/mac2mqtt/mac2mqtt.yaml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/mac2mqtt</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/mac2mqttd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/mac2mqttd.error.log</string>
</dict>
</plist>
